# root/etc/s6-overlay/s6-rc.d/fixurl/up
#!/command/execlineb -P
# export env from /etc/cont-env and friends
with-contenv

# log to stdout
fdmove -c 2 1

# run bash script
/bin/bash -e /etc/s6-overlay/scripts/fixurl.sh

