# root/etc/services.d/flask/run
#!/usr/bin/with-contenv sh

export FLASK_DEBUG=1
export FLASK_APP=app.py
export FLASK_APP=/app/app.py
exec flask run --host=0.0.0.0 --port=80
